# minecraft.service.default
[Unit]
Description=Minecraft Server
After=network.service

[Service]
User=minecraft
Group=minecraft
WorkingDirectory=/usr/local/minecraft/
ExecStart=/usr/lib/jvm/java-8-openjdk-amd64/jre/bin/java -Xms64M -Xmx1536M -jar /usr/local/minecraft/server.jar nogui
ExecStop=kill -9 $MAINPID
RestartSec=15
Restart=always

[Install]
WantedBy=multi-user.target