# minecraft.service.default
[Unit]
Description=Minecraft Server
After=network.service

[Service]
User=minecraft
Group=minecraft
Type=simple
WorkingDirectory=/usr/local/minecraft/
ExecStart=java -Xms64M -XmX1536M -ja /usr/local/minecraft/server.jar -o true
ExecStop=kill -9 $MAINPID
PIDFile=/usr/local/teamspeak/ts3server.pid
RestartSec=15
Restart=always

[Install]
WantedBy=multi-user.target