# terraria.service.default
[Unit]
Description=Terraria Server
After=network.service

[Service]
User=terraria
Group=terraria
Type=forking
KillMode=none
WorkingDirectory=/usr/local/terraria
ExecStart=/usr/bin/screen -dmS terraria /bin/bash -c "/usr/local/terraria/TerrariaServer.bin.x86_64 -config terraria.config"
ExecStop=/usr/local/bin/terrariad exit
RestartSec=15
Restart=always

[Install]
WantedBy=multi-user.target